# onfailure_alert/templates/onfailure-alert@.service
[Unit]

Description=OnFailure Alert %i

[Service]

Type=oneshot

ExecStart={{ onfailure_alert__script_path }} %i

TimeoutStartSec=5s
