# onfailure_alert/templates/onfailure-alert-example.service
[Unit]

Description=OnFailure Alert Example

StartLimitIntervalSec=120

StartLimitBurst=5

OnFailure={{ onfailure_alert__systemd_unit_name | replace('@', '@%n') }}

[Service]

Type=simple

ExecStart=/usr/bin/sh -c 'echo "It is time to fail." && exit 1'

Restart=always

RestartSec=10
